# srvrs-limadev/continuidade-operacional/systemd/limadev-heartbeat-report.timer
[Unit]
Description=LimaDev Heartbeat Report Timer

[Timer]
OnCalendar=*-*-* 07:30:00
RandomizedDelaySec=900
Persistent=true
Unit=limadev-heartbeat-report.service

[Install]
WantedBy=timers.target
